# scripts/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:web]
command=node /app/apps/web/server.js
directory=/app
user=appuser
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web-error.log
environment=NODE_ENV="production",PORT="%(ENV_PORT_WEB)s"

[program:api]
command=node /app/api-prod/dist/index.js
directory=/app/api-prod
user=appuser
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api-error.log
environment=NODE_ENV="production",PORT="%(ENV_PORT_API)s"

[program:cron]
command=/bin/bash -c 'while true; do /app/scripts/cron-reminders.sh; sleep 900; done'
directory=/app
user=appuser
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/cron.log
stderr_logfile=/var/log/supervisor/cron-error.log
